# roles/mbtvdaemon/templates/mbtvdaemon.service
[Unit]
Description=Multi-billing mbtvdaemon
After=syslog.target network.target mysql.service

[Service]
User={{ mbtvdaemon.system_user }}
RuntimeDirectory={{ mbtvdaemon.system_user }}
Type=forking
TimeoutSec=15s
PIDFile={{ mbtvdaemon.params.production.pid }}
ExecStart={{ php_path.stdout }} {{ billing.dirs.home }}/{{ mbtvdaemon.dirs.module }}/{{ mbtvdaemon.dirs.production }}/mbtvdaemon.php
ExecReload=/bin/kill -s SIGHUP $MAINPID
ExecStop=/bin/kill -s SIGTERM $MAINPID
KillMode=process

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target