# roles/mbdaemoncore/templates/mbdaemoncore.service
[Unit]
Description=Multi-billing core
After=syslog.target network.target mysql.service

[Service]
User={{ mbdaemoncore.system_user }}
RuntimeDirectory={{ mbdaemoncore.system_user }}
Type=forking
TimeoutSec=15s
PIDFile={{ mbdaemoncore.params.production.pid }}
ExecStart={{ php_path.stdout }} {{ billing.dirs.home }}/{{ mbdaemoncore.dirs.module }}/{{ mbdaemoncore.dirs.production }}/mbdaemoncore.php --noroot
ExecReload=/bin/kill -s SIGHUP $MAINPID
ExecStop=/bin/kill -s SIGTERM $MAINPID
KillMode=process

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target