# roles/mbworkerbot/templates/mbworkerbot_testing.service
[Unit]
Description=Multi-billing testing mbworkerbot
After=syslog.target network.target mysql.service

[Service]
User={{ mbworkerbot.system_user }}
RuntimeDirectory={{ mbworkerbot.system_user }}
Type=forking
TimeoutSec=15s
PIDFile={{ mbworkerbot.params.testing.pid }}
ExecStart={{ php_path.stdout }} {{ billing.dirs.home }}/{{ mbworkerbot.dirs.module }}/{{ mbworkerbot.dirs.testing }}/mbworkerbot.php --testing
ExecReload=/bin/kill -s SIGTERM $MAINPID
ExecStop=/bin/kill -s SIGTERM $MAINPID
KillMode=process

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target